# device-mlmmj-el9/src/main/resources/etc/cron.hourly/mlmmj-maint
#!/bin/bash

/usr/bin/mlmmj-maintd -F -d /var/spool/mlmmj

